# tools/pi-conf/config/autostart/hide-mouse.desktop
[Desktop Entry]
Type=Application
Name=Hide Mouse Cursor
Exec=unclutter -idle 1
NoDisplay=True
