# tools/pi-conf/config/autostart/serial-terminal.desktop
[Desktop Entry]
Type=Application
Name=Serial Terminal
Exec=lxterminal --command="screen /dev/ttyUSB0 115200"
